# myservice.sh
#!/bin/bash

#chkconfig: 2345 20 80

/usr/local/daemon/daemon

exit 0
